# isu13/mitamae/cookbooks/ipaddr/files/etc/systemd/system/isunarabe-env-ipaddr.service
[Unit]
Description=run isunarabe-env-ipaddr once
Before=mysql.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/isunarabe-env-ipaddr.sh

[Install]
WantedBy=multi-user.target
